# raspberry/.config/autostart/BancDessai.desktop
[Desktop Entry]
Type=Application
Name=Banc_Robot
Exec=python /home/pi/Banc_Robot/Autorun.py
